# mesh/deploy/mesh-dashboard.service
[Unit]
Description=Slancha-Mesh Streamlit replay dashboard
Documentation=https://github.com/SlanchaAi/slancha-mesh/blob/main/docs/SLANCHA_MESH_V0_SPEC.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# Replay JSONL the dashboard renders. install_dashboard.sh sets this to
# the rolling nightly-smoke history; override via `systemctl --user edit`.
Environment="MESH_DASHBOARD_REPLAY=%S/slancha-mesh/replay-history/latest.jsonl"
Environment="MESH_DASHBOARD_PORT=8501"

# slancha-mesh is dev-installed from source; systemd runs from / without
# carrying PYTHONPATH or cwd. Set both so streamlit can import the mesh
# package. Override via `systemctl --user edit` if the checkout is elsewhere.
WorkingDirectory=%h/Source/slancha-mesh
Environment="PYTHONPATH=%h/Source/slancha-mesh"

# Bind to 127.0.0.1 only. External exposure is handled by the
# mesh-dashboard-tunnel.service (cloudflared); LAN reach is intentionally
# blocked so the only ingress path is through Cloudflare's edge.
ExecStart=/usr/bin/env streamlit run \
    %h/Source/slancha-mesh/mesh/dashboard/streamlit_app.py \
    --server.address 127.0.0.1 \
    --server.port ${MESH_DASHBOARD_PORT} \
    --server.headless true \
    --browser.gatherUsageStats false \
    -- \
    --replay ${MESH_DASHBOARD_REPLAY}

# Streamlit responds to SIGTERM cleanly; give it 10s to flush before kill.
TimeoutStopSec=10
KillMode=mixed

# Restart on crash, but back off to avoid rapid-fire loops.
Restart=on-failure
RestartSec=15

# Sandboxing — mirrors the mesh-nightly-smoke posture (spark-v005-nightly-install-fix).
# ProtectHome=read-only would block streamlit from importing the source
# checkout; ProtectHome=tmpfs gives a private view of $HOME and the explicit
# BindReadOnlyPaths whitelists ONLY the slancha-mesh checkout.
PrivateTmp=true
ProtectSystem=full
ProtectHome=tmpfs
BindReadOnlyPaths=%h/Source/slancha-mesh
NoNewPrivileges=true

# Replay-history dir lives under XDG_STATE_HOME and needs to be readable;
# expose it via BindReadOnlyPaths since ProtectHome=tmpfs hides it by default.
BindReadOnlyPaths=%S/slancha-mesh

[Install]
WantedBy=default.target
